# install/deb/prerm
systemctl stop bitd >/dev/null 2>&1
systemctl disable bitd >/dev/null 2>&1
